# systemd/.config/systemd/user/checkupdates.service
[Unit]
Description=Check for updates

[Service]
Type=oneshot
ExecStart=/home/p00f/bin/notify_updates.sh