# systemd/.config/systemd/user/checkupdates.timer
[Unit]
Description=Update checking timer

[Timer]
OnBootSec=2m
OnUnitActiveSec=60m
Unit=checkupdates.service

[Install]
WantedBy=timers.target
